# packages/birdscarer/files/etc/uci-defaults/010-dev-disabled-services
#!/bin/sh
#
# Copyright (C) 2016 DEVELICT Solutions s.r.o.
#

# Firewall disable
/etc/init.d/firewall disable

# Dnsmasq disable
#/etc/init.d/dnsmasq disable
